# Dockerfiles/parser/etc/supervisord.d/python.ini
[program:python]
process_name = python3
command=python3 /scripts/postfix-filter-loop/postfix-filter-loop.py
startsecs=0
autorestart=false
