# Dockerfiles/ocr/etc/supervisord.d/python.ini
[program:python]
process_name = python2
command=python2 /opt/python_ocr_tutorial/flask_server/app.py
startsecs=0
autorestart=false
