# Dockerfiles/ocr/etc/supervisord.d/readlog.ini
[program:readlog]
command=/usr/bin/tail -f /opt/python_ocr_tutorial/flask_server/error.log
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
